# sparse/etc/systemd/user/maliit-server.service.d/pipa.conf
[Unit]
Requires=
After=dbus.socket
Wants=

[Service]
ExecStart=
ExecStart=/usr/bin/maliit6-server -allow-anonymous
Environment=QT_PLUGIN_PATH=/usr/lib64/qt6/plugins
Environment=QML2_IMPORT_PATH=/usr/lib64/qt6/qml
Environment=QT_QPA_PLATFORM=wayland
Restart=on-failure
RestartSec=3
# Do not Restart=always — avoids user-session thrash if Wayland is not ready yet
StartLimitIntervalSec=60
StartLimitBurst=5
