# sparse/usr/lib/systemd/system/hexagonrpcd-sdsp.service.d/20-pipa-limits.conf
#!/bin/sh
# Soften hexagonrpcd restart storms (can contribute to boot instability on pipa)
[Service]
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=120
StartLimitBurst=5
